# stack/services/backup.timer
[Unit]
Description=schedule backup service timer

[Timer]
OnCalendar=hourly

[X-Fleet]
MachineOf=backup.service
