# stack/services/count.timer
[Unit]
Description=pollster scheduler statistics
Requires=count.service

[Timer]
OnCalendar=*:0/1

[X-Fleet]
MachineOf=count.service
