# stack/services/heartbeat.timer
[Unit]
Description=pollster health checks
Requires=heartbeat.service

[Timer]
OnCalendar=*:0/1

[X-Fleet]
MachineOf=heartbeat.service
