# stack/services/pull.timer
[Unit]
Description=Docker image updater timer

[Timer]
OnCalendar=hourly

[X-Fleet]
Global=true