# stack/services/snapshot.service
[Unit]
Description=Docker image updater

[Service]
ExecStart=/usr/bin/bash /home/core/etcsnapshot

[X-Fleet]
Global=true