# stack/services/snapshot.timer
[Unit]
Description=Etcd synchronization timer

[Timer]
OnCalendar=minutely

[X-Fleet]
Global=true