# stack/services/submit.timer
[Unit]
Description=article submitter timer
Requires=submit.service

[Timer]
OnCalendar=*:0/1

[X-Fleet]
MachineOf=submit.service
