# stack/services/submit.timer
[Unit]
Description=article submitter timer
Requires=submitter.service

[Timer]
OnCalendar=*:0/1

[X-Fleet]
MachineOf=submitter.service
