# stack/services/submitter.timer
[Unit]
Description=article submitter timer
Requires=submitter.service

[Timer]
OnCalendar=*:0/5

[X-Fleet]
MachineOf=submitter.service