# stack/services/update.timer
[Unit]
Description=Docker image updater timer
Requires=update.service

[Timer]
OnCalendar=minutely

[X-Fleet]
Global=true