# stack/services/update.timer
[Unit]
Description=Docker image updater timer

[Timer]
OnCalendar=minutely

[X-Fleet]
Global=true